# packages/pkgbuild/_000_unused/anbox-helper/systemd/anbox.service
[Unit]
Description=Anbox Session Manager

[Service]
ExecStart=anbox session-manager --window-size=720,1024 --gles-driver=host

[Install]
WantedBy=default.target
